# ptm/pyptm/runscript/ptmregistry
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

NAME="ptmregistry"
DAEMON="/usr/local/sbin/$NAME"
PIDFILE="/var/run/$NAME.pid"

depend() {
	need net
}

start() {
	ebegin "Starting $NAME"
	start-stop-daemon --start --oknodo --quiet -b --make-pidfile --pidfile "$PIDFILE" -a "$DAEMON" -x $NAME
	rv=$?
	einfo $rv
	eend $rv
}

stop() {
	ebegin "Stopping $NAME"
	start-stop-daemon --stop --quiet --pidfile "$PIDFILE" 
	rv=$?
	rm -f "$PIDFILE"
	eend $rv
}

restart() {
	svc_stop;
	svc_start;
}

